# pi/systemd/flagpole.service
[Unit]
Description=Truth Social Flagpole
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/flagpole.env
WorkingDirectory=/home/teddyh/bigbeautifulflag/pi
ExecStart=/home/teddyh/bigbeautifulflag/pi/bbf-venv/bin/python /home/teddyh/bigbeautifulflag/pi/flagpole.py
Restart=on-failure
RestartSec=5
# Piomatter pokes PIO registers directly; it needs either root or the
# cap_sys_rawio capability on the python binary. Running as root is simplest.
User=root

[Install]
WantedBy=multi-user.target
